# infra/ansible/templates/polkadot-alice.service.j2
[Unit]
Description=Polkadot Node

[Service]
ExecStart=/opt/artemis/bin/polkadot \
         --name {{ nodes.alice.name }} \
         --chain /opt/artemis/rococo.json \
         --base-path /opt/artemis/data \
         --rpc-cors=all \
         --alice \
         --node-key {{ nodes.alice.node_key }}

Restart=always

[Install]
WantedBy=multi-user.target
